# packages/pes-sixad/97-sixpair.rules
SUBSYSTEM=="usb", ATTRS{idVendor}=="054c", ENV{idProduct}="0268", RUN+="/opt/pes/sbin/sixpair"
SUBSYSTEM=="input", ATTR{name}=="PLAYSTATION(R)3 Controller*", RUN+="/opt/pes/sbin/sixpair"
